# services/nico.timer
[Unit]
Description=Nico timer

[Timer]
OnUnitActiveSec=8h
OnBootSec=1h

[Install]
WantedBy=timers.target