# services/nico.service
[Unit]
Description=nico ...

[Service]
Type=oneshot
# update this path
ExecStart=python main.py
